# tests/util/sandbox/egress_guard_e2e/multi_service-supervisord.conf
[supervisord]
nodaemon=true

[program:internal-http]
command=/usr/bin/python3 -m http.server 8080 --bind 127.0.0.1 --directory /srv/internal
user=service
startsecs=0
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
